# wiki/supervisord.conf
[supervisord]
nodaemon=true

[program:wiki]
command=/confluence/confluence/bin/start-confluence.sh -fg
